# infra/production/monitoring/systemd/alertmanager-discord.service
[Unit]
Description=Alertmanager Discord Webhook Bridge
After=network-online.target alertmanager.service
Wants=network-online.target

[Service]
Type=simple
User=alertmanager
Group=alertmanager
EnvironmentFile=/opt/frontman/monitoring/alertmanager-discord.env
ExecStart=/usr/local/bin/alertmanager-discord
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
